# build/systemd/excubitor.service
[Unit]
Description=Service to start and stop Excubitor
SuccessAction=